# templates/systemd-pull.timer
[Unit]
Description=claude-code-memory pull every {{INTERVAL}}s

[Timer]
OnBootSec=30
OnUnitActiveSec={{INTERVAL}}
AccuracySec=5
Unit=claude-code-memory-pull.service

[Install]
WantedBy=timers.target
